# hackchat/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx

[program:chat-server]
command=node server.js

[program:chat-client]
command=http-server /hack.chat/client/
